# services/backup.timer
[Unit]
Description=Run a daily backup of all VMs.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
